# files/etc/uci-defaults/99-remote-agent-bootstrap
#!/bin/sh

set -e

readonly raFunctions='/usr/lib/remote-agent/functions.sh'
if [ -s "$raFunctions" ]; then
# shellcheck source=../../usr/lib/remote-agent/functions.sh
	. "$raFunctions"
else
	logger -t 'remote-agent' 'error' "Remote-Agent Functions file ($raFunctions) not found!"
fi

mkdir -p "$STATE_DIR"; chmod 700 "$STATE_DIR"

# system user (no shell, no home)
id "$RPC_USER" >/dev/null 2>&1 || adduser -D -H -s /bin/false "$RPC_USER" 2>/dev/null || true

# Bootstrap password for initial contact (will be rotated/removed upon adoption)
PASS="$(dd if=/dev/urandom bs=24 count=1 2>/dev/null | base64 | tr -d '\n=/+' | cut -c1-32)"
HASH="$(openssl passwd -6 "$PASS")"

uci add rpcd login >/dev/null
uci set rpcd.@login[-1].username="$RPC_USER"
uci set rpcd.@login[-1].password="$HASH"
uci add_list rpcd.@login[-1].read="remote-agent"
uci add_list rpcd.@login[-1].write="remote-agent"
uci commit rpcd

cat >"$CRED_JSON" <<EOF
{"username":"$RPC_USER","password":"$PASS"}
EOF
chmod 600 "$CRED_JSON"

# Ensure /ubus is exposed
uci -q show uhttpd | grep -q "ubus_prefix='/ubus'" || {
	uci add_list uhttpd.main.ubus_prefix='/ubus'
	uci commit uhttpd
	/etc/init.d/uhttpd restart || true
}

/etc/init.d/rpcd restart >/dev/null 2>&1 || true

# Pre-generate an adoption token used in discoverable mode (one-time)
if [ ! -s "$ADOPT_TOKEN_FILE" ]; then
	dd if=/dev/urandom bs=18 count=1 2>/dev/null | base64 | tr -d '\n=/' | cut -c1-24 > "$ADOPT_TOKEN_FILE"
	chmod 600 "$ADOPT_TOKEN_FILE"
fi

exit 0
